# systemd/model-worker.service
[Unit]
Description=DevOps Assignment - Model Worker
After=network.target

[Service]
Type=simple
User=devops
Group=devops
EnvironmentFile=/etc/devops-assignment.env
WorkingDirectory=/opt/devops-assignment/app/workers
ExecStart=/opt/devops-model-venv/bin/uvicorn model_worker:app --host 0.0.0.0 --port ${MODEL_WORKER_PORT}
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
